# systemd/robot.service
[Unit]
Description=Robot Controller Service
After=network.target
[Service]
ExecStart=/usr/bin/python3 /path/to/robot_controller/app/main.py
Restart=always
[Install]
WantedBy=multi-user.target